# led.service
[Unit]
Description=Affiche led rouge et led verte

[Service]
Type=simple
ExecStart=/home/pi/svxlink/led.sh


[Install]
WantedBy=multi-user.target
